# pluto-backend/deploy/systemd/pluto-migrator.service
[Unit]
Description=Pluto Backend — Database Migration Runner (one-shot)
Requires=docker.service
After=docker.service network-online.target pluto-backend.service
Wants=network-online.target

# একবার চলে exit করবে; ব্যর্থ হলে 30s পর ৩ বার retry।
[Service]
Type=oneshot
WorkingDirectory=/root/backend-joy/pluto-backend
# docker-compose.migrator.yml একটি overlay যা `migrator` service কে চালায়
# (packages/api/scripts/migrate.mjs) এবং শেষ হলে exit করে।
ExecStart=/bin/bash deploy/preflight-migrations.sh

# ব্যর্থ হলে স্বয়ংক্রিয় retry (systemd oneshot এ Restart=on-failure সমর্থিত নয়,
# তাই আমরা pluto-migrator.timer দিয়ে periodic retry করি — নিচে দেখুন)।
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pluto-migrator
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
